# conf/desktop/etc/profile.d/couldinho-desktop.sh
#!/bin/sh

# The system being installed
export COULDINHO_SYSTEM="desktop"
