# pkg/couldinho-vmware.install
post_install() {
    post_upgrade
}

post_upgrade() {

    # System environment variables
    system=vmware
    if [[ "${COULDINHO_SYSTEM}" != *"${system}"* ]]; then
        export "COULDINHO_SYSTEM=${system},${COULDINHO_SYSTEM}" \
            > /etc/environment.d/80-couldinho-system
    fi

    # Enable services
    systemctl enable vmtoolsd
    systemctl enable vmware-vmblock-fuse
}

