# pkg/couldinho-vmware.install
post_install() {
    post_upgrade
}

post_upgrade() {

    # Enable services
    systemctl enable vmtoolsd
    systemctl enable vmware-vmblock-fuse
}

