# pkg/couldinho-dev.install
post_install() {
    post_upgrade
}

post_upgrade() {

    # System environment variables
    system=dev
    if [[ "${COULDINHO_SYSTEM}" != "${system}" ]]; then
        export "COULDINHO_SYSTEM=${system}" \
            > /etc/environment.d/80-couldinho-system
    fi

    # Enable services
    systemctl enable docker.service
}
