# pkg/couldinho-dev.install
post_install() {
    post_upgrade
}

post_upgrade() {

    # System environment variables
    mv /etc/profile.d/couldinho-desktop.sh /etc/profile.d/couldinho.sh

    # Enable services
    systemctl enable docker.service
}
