# pkg/couldinho-vmware.install
post_install() {
    post_upgrade
}

post_upgrade() {

    # Ensure vmware-user gets started
    cp /etc/i3/couldinho-config /etc/i3/config

    # Enable services
    systemctl enable vmtoolsd
    systemctl enable vmware-vmblock-fuse
}

