# conf/desktop/etc/systemd/user/waybar-updates.service
[Unit]
Description=Refresh list of outdated packages in waybar

[Service]
ExecStart=/usr/local/bin/waybar-updates refresh
